# scripts/launchd/com.gdist.vault-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.gdist.vault-sync.plist
  Synker Obsidian-vault til GitHub hvert 10. minutt.

  Installer:
    cp ~/Documents/GlobalDistribution/scripts/launchd/com.gdist.vault-sync.plist \
       ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.gdist.vault-sync.plist

  Status:
    launchctl list | grep gdist

  Logg:
    tail -f /tmp/gdist-vault-sync.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gdist.vault-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/gdist/Documents/GlobalDistribution/scripts/sync-all.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/gdist-vault-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/gdist-vault-sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/gdist</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
